# src/certferry-renew.service
[Unit]
Description=cert-ferry certificate renewal

[Service]
Type=oneshot
ExecStart=%CERTFERRY_EXE% --renew
